# ros3d-kr.conf
description	"Ros3D KR"

start on started dbus

respawn

console log
exec /usr/bin/ros3d-kr --http-port 8090
